# container/crontab
* * * * * /etc/cloudwatch/monitor-resources.sh
